# deployment/instance/documatcher.timer
[Unit]
Description=Keep DocuMatcher up to date with main

[Timer]
OnBootSec=0s
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=timers.target
